# .config/systemd/user/mpd-notification-daemon.service
[Unit]
Description=Simple MPD notification daemon
After=mpd.service

[Service]
ExecStart=%h/.local/bin/mpdnotdaem

[Install]
RequiredBy=mpd.service
