# .config/systemd/user/mpdscribble.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpdscribble --conf %h/.config/mpdscribble/mpdscribble.conf --no-daemon

[Install]
WantedBy=
RequiredBy=mpd.service
