# .config/systemd/user/mpdscribble.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpdscribble --conf /home/nicolai/.config/mpdscribble/mpdscribble.conf --no-daemon